# Config_files/servicios_y_scripts_toggle_tva_games/etc/systemd/system/emulationstation-session.service
[Unit]
Description=EmulationStation / RetroPie session
After=graphical.target

[Service]
Type=simple
# IMPORTANT: Replace <TVARGENTA_USER> with the actual username (e.g., 'pi' on standard Raspberry Pi)
User=<TVARGENTA_USER>
Group=<TVARGENTA_USER>
Environment=HOME=/home/<TVARGENTA_USER>
Environment=USER=<TVARGENTA_USER>
Environment=LOGNAME=<TVARGENTA_USER>
Environment=SHELL=/bin/bash
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=TERM=xterm
ExecStart=/opt/retropie/supplementary/emulationstation/emulationstation
Restart=on-abort

[Install]
WantedBy=multi-user.target

