# Config_files/servicios_y_scripts_toggle_tva_games/etc/systemd/system/tvargenta.service
# /etc/systemd/system/tvargenta.service
[Unit]
Description=TVArgenta Backend (app.py)
Wants=network-online.target
After=network-online.target

[Service]
# IMPORTANT: This is a template file for reference only.
# The install.sh script generates this service dynamically with the correct user.
# Do not use this file directly; run install.sh for proper setup.
User=<TVARGENTA_USER>
WorkingDirectory=/srv/tvargenta
ExecStart=/srv/tvargenta/venv/bin/python /srv/tvargenta/app.py

# PATH completo para que encuentre pkill y demás
Environment="PATH=/srv/tvargenta/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="FLASK_ENV=production"

Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal
KillSignal=SIGINT
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target

